# systemd/.config/systemd/user/backup.service
[Unit]
Description=Backup
Wants=backup.timer

[Service]
Type=oneshot
ExecStart=env -i HOME=/home/user /home/user/.local/bin/borgmatic --verbosity 2
